# data/service/thm.service
[Unit]
Description=thm service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/python /.data/service/server.py
Restart=always
WorkingDirectory=/.data/service/
Nice=19
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target